# com.missioncontrol.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.missioncontrol.server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>--enable-source-maps</string>
        <string>--max-old-space-size=256</string>
        <string>/Users/elf-6/mission-control-server/server.js</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <true/>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/elf-6/mission-control-server/logs/server.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/elf-6/mission-control-server/logs/server.error.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/elf-6/mission-control-server</string>
    <key>ProcessType</key>
    <string>Background</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>
</dict>
</plist>